# gitlab-ci/core_service.conf
[program:ubusd]
priority=1
command=/bin/bash -c "/usr/sbin/ubusd"

[program:rpcd]
priority=2
command=/bin/bash -c "/usr/sbin/rpcd"
